# installer/qlcplus-podman.service
[Unit]
Description=QLC+ VNC Openbox Container (Podman)
After=network.target systemd-udev-settle.service
Requires=systemd-udev-settle.service
Before=graphical.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=josh
WorkingDirectory=/opt/StageKitPied/podman

# Ensure device is ready before launch
ExecStartPre=/usr/bin/udevadm settle --timeout=60
ExecStartPre=/bin/bash -c 'for i in {1..10}; do [ -e /dev/ftdi-A10JGB3H ] && exit 0; sleep 2; done; echo "FTDI device not found" >&2; exit 1'

# Ensure device is released after stoppage
ExecStopPost=/opt/StageKitPied/udev/resetFtdiDmx.sh

# Start and stop scripts
ExecStart=/opt/StageKitPied/podman/run_podman.sh
ExecStop=/opt/StageKitPied/podman/destroy_podman.sh

# Logging
StandardOutput=append:/var/log/qlcplus-podman.log
StandardError=inherit

# Restart policy: don’t automatically restart unless manually triggered
Restart=no

[Install]
WantedBy=multi-user.target

